# getemail.service
[Unit]
Description=get email

[Service]
ExecStart=/home/tovare/getemail/linux-getemail

[Install]
WantedBy=multi-user.target